# board/dacha/rootfs_overlay/etc/udev/rules.d/10-usb-nic.rules
ACTION=="add", KERNEL=="eth1", SUBSYSTEM=="net", RUN+="/sbin/ifup eth1"